# rocky8/sshd.ini
[program:sshd]
command=/usr/sbin/sshd -D
directory=/bin
user=root
autostart=true
autorestart=false